# tg-archive-bot.service
[Unit]
Description=Telegram Archive Bot
After=network.target

[Service]
User=archivebot
WorkingDirectory=/opt/tg-archive-bot
ExecStart=/opt/tg-archive-bot/venv/bin/python bot.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
